# extracted_fs/etc/init.d/atcmd
#!/bin/sh /etc/rc.common

STOP=90

stop(){
	/usr/bin/atcmd
}
