# extracted_fs/usr/lib/opkg/info/lte-telephony.list
/usr/include/paths_defs.h
/usr/include/hawk_api.h
/usr/include/shm_share.h
/usr/include/linux_types.h
/usr/include/global_types.h
/usr/include/ml_utils.h
/usr/include/ci_stub_ttc_macro.h
/usr/include/NVMServer_api.h
/usr/include/ppprd_linux.h
/usr/include/nvm_shared.h
/usr/lib/libmtel.so
/etc/init.d/lte-telephony-prod
/usr/include/diag_ioctl.h
/etc/telinit
/etc/init.d/lte-telephony
/usr/include/common_datastub_macro.h
/usr/include/eeh_ioctl.h
/usr/include/linux_driver_types.h
/usr/include/fdi_ram.h
/usr/include/pxa9xx_cp.h
/usr/include/gbl_types.h
/usr/include/pxa_cp_load_ioctl.h
